# puppet/modules/oracle_xe/files/etc/profile.d/oracle-xe.sh

# Source Oracle environment
if [ -r /u01/app/oracle/product/*/xe/bin/oracle_env.sh ]
then
    . /u01/app/oracle/product/*/xe/bin/oracle_env.sh
fi

